# gunicorn.service
[Unit]
Description=PM Knowledge Hub Gunicorn Service
After=network.target

[Service]
User=root
Group=root
WorkingDirectory=/root/.openclaw/workspace/projects/pm-knowledge-hub
ExecStart=/root/.openclaw/workspace/projects/pm-knowledge-hub/venv/bin/gunicorn \
    --access-logfile - \
    --workers 4 \
    --bind 127.0.0.1:8000 \
    config.wsgi:application

[Install]
WantedBy=multi-user.target
